# mgmt/etc/supervisor.d/vcld.conf
[supervisord]
nodaemon=true

[program:vcld]
command=/usr/local/vcl/bin/vcld -vvv -conf=/etc/vcl/vcld.conf
