# mgmt/etc/supervisor.d/vcld.conf
[supervisord]
nodaemon=true

[program:vcld]
command=/usr/local/vcl/bin/vcld --verbose  --config=/etc/vcl/vcld.conf --nodaemon
