# www/etc/supervisor.d/postfix.conf
[supervisord]
nodaemon=true

[program:postfix]
process_name = postfix
directory = /etc/postfix
command	= /postfix-wrapper.sh
startsecs = 0
autorestart = false